# src/papa_bringup/etc/ros-system-bringup.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/ros-system-bringup.sh

[Install]
WantedBy=default.target
